# deploy/sge-traefik.service
[Unit]
Description=SGE — Traefik reverse proxy
Documentation=https://traefik.io
After=network.target
Before=sge.service sge-frontend.service

[Service]
Type=simple
User=traefik
Group=sge
# CAP_NET_BIND_SERVICE allows binding :80/:443 without root
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

WorkingDirectory=/etc/sge/traefik

ExecStart=/usr/local/bin/traefik --configFile=/etc/sge/traefik/traefik.yml
ExecReload=/bin/kill -HUP $MAINPID

Restart=on-failure
RestartSec=5s
StartLimitInterval=60s
StartLimitBurst=3

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/etc/sge/traefik /var/log/sge

# Logging (FHS: logs en /var/log/sge/)
StandardOutput=append:/var/log/sge/traefik-stdout.log
StandardError=append:/var/log/sge/traefik-error.log

[Install]
WantedBy=multi-user.target
